# server/etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#

su redline -c 'screen -d -m -S screenslime -c /home/redline/webapps/screenrc'
